# mac_env/startup/conf.d/sshd.conf
[program:sshd]
command=/usr/sbin/sshd -D -p 7777
autorestart=true