# mac_env/startup/conf.d/sshd.conf
[program:sshd]
command=/usr/sbin/sshd -D -p %(ENV_GDB_SSH_PORT)s
autorestart=true